# share/udev/rules.d/00-keymap.rules
ACTION=="add" SUBSYSTEM=="input" ATTRS{name}=="Lite-On Technology Corp. USB Keyboard" RUN+="keymap $name lenovo-K4803"
