# etc/_config/_systemd/_user/kupfer.service
[Service]
Type = simple
ExecStart = /usr/bin/kupfer --no-splash
Restart = on-failure

[Install]
WantedBy = default.target
