# share/systemd/system/dnscrypt-proxy.service.erb
[Unit]
Description=Encrypted DNS proxy daemon
Requires=loopback-alias.service
After=loopback-alias.service

[Service]
Type=simple
ExecStart=<%= path 'dnscrypt-proxy' %> -a 127.0.0.2 -r 208.67.220.220 -u dns

[Install]
WantedBy=multi-user.target
