# share/systemd/system/extra-loopback-addresses.service.erb
[Unit]
Description=Add extra loopback addresses
Before=network.target

[Service]
Type=oneshot
ExecStart=<%= path 'ip' %> addr add 127.0.0.53/8 dev lo
ExecStop=<%= path 'ip' %> addr del 127.0.0.53/8 dev lo
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
