# share/systemd/system/hardware-settings.service.erb
[Unit]
Description=Hardware settings

[Service]
Type=oneshot
ExecStart=/etc/local/hardware-settings.sh

[Install]
WantedBy=multi-user.target
