# share/systemd/system/powersave.service.erb
[Unit]
Description=Power saving features

[Service]
Type=oneshot
ExecStart=/etc/local/powersave.sh

[Install]
WantedBy=multi-user.target
