# share/systemd/system/setuptmpfiles.service.erb
[Unit]
Description=Move sensitive user files to temporary directories

[Service]
Type=oneshot
ExecStart=<%= path 'setuptmpfiles' %> -u root,guns
ExecStop=<%= path 'setuptmpfiles' %> -u root,guns
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
