# share/systemd/system/sleep-gpg-agent.service.erb
[Unit]
Description=Send SIGHUP to gpg-agent processes on suspend

[Service]
Type=oneshot
ExecStart=<%= path 'killall' %> -HUP gpg-agent

[Install]
WantedBy=sleep.target
