# share/systemd/system/sleep-gpg-agent.service.erb
[Unit]
Description=Send SIGHUP to gpg-agent processes on suspend
Before=sleep.target

[Service]
Type=oneshot
ExecStart=<%= path 'pkill' %> -x -HUP gpg-agent

[Install]
WantedBy=sleep.target
