# share/systemd/system/sleep-sysrq-mask.service.erb
[Unit]
Description=Set magic SysRq command mask on suspend
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/cat /proc/sys/kernel/sysrq > /run/sysrq.state; echo 128 > /proc/sys/kernel/sysrq'
ExecStop=/bin/sh -c '/usr/bin/cat /run/sysrq.state > /proc/sys/kernel/sysrq; /usr/bin/rm -f /run/sysrq.state'
RemainAfterExit=yes

[Install]
WantedBy=sleep.target
