# share/systemd/system/sleep-slock.service.erb
[Unit]
Description=Screen lock on sleep
Before=sleep.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=<%= path 'slock' %>

[Install]
WantedBy=sleep.target
